# cron/crontab
AWS_CONFIG_FILE="/home/ubuntu/.aws/config"
*/15 * * * * /home/ubuntu/cron/sync-blocks.sh
